# upstart/memcached.conf
description "memcached"

env MEMCACHED=/usr/bin/memcached
env host=xxx

start on runlevel [2345]
# Not sure why it was recommended to use ^ rather than !.  I'm sticking with !.
stop on runlevel [!2345]

# This test is completely optional, I'm just paranoid.
pre-start script
  test -x $MEMCACHED || { stop; exit 0; }
end script

respawn
exec $MEMCACHED -m 64 -p 11211 -u memcache -l $host
